# oracle/files/rc.local
#!/bin/bash

clear_user(){
	for user in oracle; do
		userdel --force --remove "$user" || passwd -l "$user" || true
		rm -rf "/home/$user/"
		rm -f "/etc/sudoers.d/oracle"
	done
}

clean_up(){
	rm -rf /etc/rc.local
	rm -rf /etc/rc.d/rc.local

	systemctl disable rc-local
	systemctl stop rc-local
	rm -rf /etc/systemd/system/rc-local.service
	rm -rf /etc/subgid.lock
}

clear_user

clean_up
